# WakeMyPotata/src/wmp.timer
# timeout 600 (Do not Modify this auto-generated line!)
# threshold 10 (Do not Modify this auto-generated line!)

[Unit]
Description=WakeMyPotata: Emergency RAID/system safe shutdown (battery-only)

[Timer]
Unit=wmp.service
OnBootSec=0
OnCalendar=*-*-* *:00/01:00
Persistent=true

[Install]
WantedBy=timers.target
