# WakeMyPotata/src/wmp.service
[Unit]
Description=WakeMyPotata

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/wmp-run --timeout 600
