# WakeMyPotata/src/wmp.service
[Unit]
Description=WakeMyPotata: Service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/wmp-run --timeout 600 && exit 0
