# WakeMyPotata/src/wmp.service
# timeout 600 (Do not Modify this auto-generated line!)
# threshold 10 (Do not Modify this auto-generated line!)

[Unit]
Description=WakeMyPotata: Emergency RAID/system safe shutdown (battery-only)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/wmp-run 600
