# WakeMyPotata/src/wmp.service
[Unit]
Description=WakeMyPotata: Emergency RAID/system safe shutdown (battery-only)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/wmp-run --threshold 10 --timeout 600
